# templates/certbot.timer.j2
[Unit]
Description=Run certbot weekly

[Timer]
OnCalendar={{ certbot_renewal_interval }}
Persistent=true

[Install]
WantedBy=timers.target
